# debian/mqrestt.service
[Unit]
Description=MQTT-REST bridge

[Service]
Type=simple
ExecStart=/usr/bin/mqrestt -c /etc/mqrestt.conf

[Install]
WantedBy=multi-user.target
